# build/openwrt_amlogic/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/sbin/balethirq.pl

## 取消自动挂载非必要分区，避免更新时出现分区识别问题
# umount /mnt/mmcblk2p1 && mount /dev/mmcblk2p1 /boot
# umount /mnt/mmcblk2p2 && umount /mnt/mmcblk2p3

exit 0
